# template.service.in
[Unit]
Description=Client application for Motorcortex
After=ethercat.target motorcortex.service
PartOf=motorcortex.service
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStartPre=/bin/sleep 3
ExecStart=@EXEC_START@
WorkingDirectory=/home/admin
User=admin
Restart=always
RestartSec=1

[Install]
WantedBy=motorcortex.service
